# Отказоустойчивость/backup/screenshots/crontab
30 2 * * * /home/admin/netlogy/homework/Отказоустойчивость/backup/backup_home.sh
